# python-image.service
[Unit]
Description=Python Image Generator Flask Application
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/python-image
Environment="PATH=/opt/python-image/venv/bin"
Environment="PORT=5566"
ExecStart=/opt/python-image/venv/bin/gunicorn -w 4 -b 0.0.0.0:5566 app:app
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/python-image/access.log
StandardError=append:/var/log/python-image/error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

